# erp-server/agents/bazooka/com.evertrust.bazooka.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.evertrust.bazooka</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/kobewannkenobi/marketing-agent-workflows/bazooka/run_bazooka.sh</string>
    </array>

    <!-- Fires at both candidate local times; the script's Berlin-hour guard
         picks whichever one is 08:00 Europe/Berlin (handles German DST). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/kobewannkenobi/marketing-agent-workflows/bazooka/runs/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kobewannkenobi/marketing-agent-workflows/bazooka/runs/launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
